# files/flatpak/etc/systemd/user/myos-flatpak-setup.service
[Unit]
Description=One-time per-user Flatpak setup
ConditionPathExists=!%h/.config/myos/flatpak-setup-done

[Service]
Type=oneshot
ExecStart=/usr/libexec/myos-flatpak-setup.sh

[Install]
WantedBy=default.target
