# files/wl/etc/systemd/system/wl-akmods-load.service
[Unit]
Description=Build & load Broadcom wl once per new kernel (ostree-safe)
After=local-fs.target
RequiresMountsFor=/usr/lib/modules/%v
ConditionPathExists=!/var/lib/wl-akmods/done-%v

[Service]
Type=oneshot
TimeoutStartSec=900
ExecCondition=/usr/local/sbin/needs-wl.sh
ExecStart=/usr/local/sbin/wl-akmods-build.sh

[Install]
WantedBy=multi-user.target
