# files/flatpak/etc/systemd/system/myos-hide-org-system.service
# /usr/lib/systemd/system/myos-hide-org-system.service
[Unit]
Description=Hide the org-system Flatpak remote from enumeration
ConditionPathExists=!/var/lib/myos/org-system-hidden
After=multi-user.target

[Service]
Type=oneshot
# Only hide if the remote exists; otherwise exit 0
ExecStart=/usr/bin/bash -lc '/usr/bin/flatpak --system remotes --columns=name | /usr/bin/grep -qx org-system && /usr/bin/flatpak --system remote-modify --no-enumerate org-system || true'
ExecStart=/usr/bin/mkdir -p /var/lib/myos
ExecStart=/usr/bin/touch /var/lib/myos/org-system-hidden

[Install]
WantedBy=multi-user.target
