# files/flatpak/etc/systemd/system/myos-hide-org-system.service
[Unit]
Description=Hide the org-system Flatpak remote from enumeration
After=network-online.target
ConditionPathExists=!/var/lib/myos/org-system-hidden

[Service]
Type=oneshot
ExecStart=/usr/libexec/myos-hide-org-system.sh

[Install]
WantedBy=multi-user.target